# config/templates/mcweb-worker.service
[Unit]
Description=McWeb Sidekiq Worker
After=network.target postgresql.service mcweb-web.service

[Service]
Type=simple
User=mcweb
WorkingDirectory=/opt/mcweb/current
EnvironmentFile=/etc/mcweb/mcweb.env
ExecStart=/home/mcweb/.rbenv/shims/bundle exec sidekiq -C config/sidekiq.yml
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
